# infra/coliseum.service
[Unit]
Description=Coliseum Autonomous Trading Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/coliseum/backend
EnvironmentFile=/home/ubuntu/coliseum/backend/.env
ExecStart=/home/ubuntu/coliseum/backend/venv/bin/python -m coliseum daemon --host 127.0.0.1 --port 9000
Restart=on-failure
RestartSec=60
TimeoutStopSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coliseum

[Install]
WantedBy=multi-user.target
